# gunicorn/gunicorn.conf
[program:gunicorn]
directory=/home/ubuntu/locallibrary
command=/home/ubuntu/env/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/locallibrary/app.sock locallibrary.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log

[group:guni]
programs:gunicorn